# fed-overlay/etc/systemd/system/pfa.slice
[Unit]
Description=PFA Experiments
DefaultDependencies=no
Before=slices.target pfa.service

[Slice]
MemoryLimit=1G
