# scripts/bootstrap/pfa.service
[Unit]
Description=PFA related initialization (e.g. setting up swap)

[Service]
Type=oneshot
ExecStart=/root/pfa-exp/scripts/fedora_init.sh

[Install]
WantedBy=default.target
